# docker/lite/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
logfile_maxbytes=50MB
logfile_backups=3
loglevel=info

; =============================================================================
; Redis Server
; In-memory data store for transcription streams and bot coordination
; =============================================================================
[program:redis]
command=/usr/bin/redis-server --bind 0.0.0.0 --port 6379 --dir /var/lib/redis --appendonly yes --appendfsync everysec
autostart=true
autorestart=true
priority=5
startsecs=2
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; =============================================================================
; Xvfb - Virtual X11 Display
; Required for Playwright browsers in headless mode
; =============================================================================
[program:xvfb]
command=Xvfb :99 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
autostart=true
autorestart=true
priority=10
startsecs=2
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; =============================================================================
; PulseAudio - Audio server for TTS playback
; Required for paplay to pipe synthesized speech into meetings
; =============================================================================
[program:pulseaudio]
command=pulseaudio --system --disallow-exit --disallow-module-loading=false --log-level=notice
autostart=true
autorestart=true
priority=12
startsecs=2
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; =============================================================================
; PulseAudio Sink Setup - creates null sinks and virtual mic
; Runs once after PulseAudio is ready
; =============================================================================
[program:pulseaudio-setup]
command=/usr/local/bin/setup-pulseaudio-sinks.sh
autostart=true
autorestart=false
startsecs=0
priority=13
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; =============================================================================
; WhisperLive - Real-time Transcription Service
; CPU-based Whisper inference server
; =============================================================================
[program:whisperlive]
command=sh -c 'if [ "%(ENV_WHISPER_BACKEND)s" = "remote" ]; then python3 run_server.py --port 9090 --backend remote; else python3 run_server.py --port 9090 --backend faster_whisper --faster_whisper_custom_model_path %(ENV_WHISPER_MODEL_SIZE)s; fi'
directory=/app/whisperlive
autostart=true
autorestart=true
priority=20
startsecs=30
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",DEVICE_TYPE="%(ENV_DEVICE_TYPE)s",CONSUL_ENABLE="false",REDIS_URL="%(ENV_REDIS_URL)s",REDIS_STREAM_URL="%(ENV_REDIS_URL)s",REDIS_STREAM_NAME="transcription_segments",LOG_LEVEL="%(ENV_LOG_LEVEL)s",TRANSCRIBER_URL="%(ENV_TRANSCRIBER_URL)s",TRANSCRIBER_API_KEY="%(ENV_TRANSCRIBER_API_KEY)s",WL_RECORDING_UPLOAD_URL="http://localhost:8080/internal/recordings/upload",WL_RECORDING_DIR="/var/lib/vexa/recordings/spool"

; =============================================================================
; Transcription Collector
; Collects transcription segments from Redis and stores in PostgreSQL
; =============================================================================
[program:transcription-collector]
command=uvicorn main:app --host 0.0.0.0 --port 8123 --log-level %(ENV_LOG_LEVEL)s
directory=/app/transcription-collector
autostart=true
autorestart=true
priority=30
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app/shared-models",REDIS_URL="%(ENV_REDIS_URL)s",REDIS_HOST="%(ENV_REDIS_HOST)s",REDIS_PORT="%(ENV_REDIS_PORT)s",REDIS_PASSWORD="%(ENV_REDIS_PASSWORD)s",DATABASE_URL="%(ENV_DATABASE_URL)s",DB_HOST="%(ENV_DB_HOST)s",DB_PORT="%(ENV_DB_PORT)s",DB_NAME="%(ENV_DB_NAME)s",DB_USER="%(ENV_DB_USER)s",DB_PASSWORD="%(ENV_DB_PASSWORD)s",REDIS_STREAM_NAME="transcription_segments",REDIS_CONSUMER_GROUP="collector_group",STORAGE_BACKEND="%(ENV_STORAGE_BACKEND)s",LOCAL_STORAGE_DIR="%(ENV_LOCAL_STORAGE_DIR)s",LOCAL_STORAGE_FSYNC="%(ENV_LOCAL_STORAGE_FSYNC)s",MINIO_ENDPOINT="%(ENV_MINIO_ENDPOINT)s",MINIO_ACCESS_KEY="%(ENV_MINIO_ACCESS_KEY)s",MINIO_SECRET_KEY="%(ENV_MINIO_SECRET_KEY)s",MINIO_BUCKET="%(ENV_MINIO_BUCKET)s",MINIO_SECURE="%(ENV_MINIO_SECURE)s",AWS_REGION="%(ENV_AWS_REGION)s",AWS_ACCESS_KEY_ID="%(ENV_AWS_ACCESS_KEY_ID)s",AWS_SECRET_ACCESS_KEY="%(ENV_AWS_SECRET_ACCESS_KEY)s",S3_BUCKET="%(ENV_S3_BUCKET)s",S3_ENDPOINT="%(ENV_S3_ENDPOINT)s",S3_SECURE="%(ENV_S3_SECURE)s"

; =============================================================================
; Admin API
; User and API token management
; Exposed on port 8057 (external) and 8001 (internal for API Gateway)
; =============================================================================
[program:admin-api]
command=uvicorn app.main:app --host 0.0.0.0 --port 8057 --log-level %(ENV_LOG_LEVEL)s
directory=/app/admin-api
autostart=true
autorestart=true
priority=40
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app/shared-models",DATABASE_URL="%(ENV_DATABASE_URL)s",DB_HOST="%(ENV_DB_HOST)s",DB_PORT="%(ENV_DB_PORT)s",DB_NAME="%(ENV_DB_NAME)s",DB_USER="%(ENV_DB_USER)s",DB_PASSWORD="%(ENV_DB_PASSWORD)s",ADMIN_API_TOKEN="%(ENV_ADMIN_API_TOKEN)s"

; =============================================================================
; Bot Manager
; Manages bot lifecycle using process orchestrator (no Docker socket needed)
; =============================================================================
[program:bot-manager]
command=uvicorn app.main:app --host 0.0.0.0 --port 8080 --log-level %(ENV_LOG_LEVEL)s
directory=/app/bot-manager
autostart=true
autorestart=true
priority=50
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app/shared-models",ORCHESTRATOR="process",BOT_SCRIPT_PATH="/app/vexa-bot/dist/docker.js",BOT_WORKING_DIR="/app/vexa-bot",WHISPER_LIVE_URL="ws://localhost:9090",BOT_CALLBACK_BASE_URL="http://localhost:8080",DISPLAY=":99",REDIS_URL="%(ENV_REDIS_URL)s",DATABASE_URL="%(ENV_DATABASE_URL)s",DB_HOST="%(ENV_DB_HOST)s",DB_PORT="%(ENV_DB_PORT)s",DB_NAME="%(ENV_DB_NAME)s",DB_USER="%(ENV_DB_USER)s",DB_PASSWORD="%(ENV_DB_PASSWORD)s",ADMIN_TOKEN="%(ENV_ADMIN_API_TOKEN)s",TRANSCRIPTION_COLLECTOR_URL="http://localhost:8123",TTS_SERVICE_URL="http://localhost:8059",STORAGE_BACKEND="%(ENV_STORAGE_BACKEND)s",LOCAL_STORAGE_DIR="%(ENV_LOCAL_STORAGE_DIR)s",LOCAL_STORAGE_FSYNC="%(ENV_LOCAL_STORAGE_FSYNC)s",MINIO_ENDPOINT="%(ENV_MINIO_ENDPOINT)s",MINIO_ACCESS_KEY="%(ENV_MINIO_ACCESS_KEY)s",MINIO_SECRET_KEY="%(ENV_MINIO_SECRET_KEY)s",MINIO_BUCKET="%(ENV_MINIO_BUCKET)s",MINIO_SECURE="%(ENV_MINIO_SECURE)s",AWS_REGION="%(ENV_AWS_REGION)s",AWS_ACCESS_KEY_ID="%(ENV_AWS_ACCESS_KEY_ID)s",AWS_SECRET_ACCESS_KEY="%(ENV_AWS_SECRET_ACCESS_KEY)s",S3_BUCKET="%(ENV_S3_BUCKET)s",S3_ENDPOINT="%(ENV_S3_ENDPOINT)s",S3_SECURE="%(ENV_S3_SECURE)s"

; =============================================================================
; API Gateway
; Main entry point - routes requests to appropriate services
; =============================================================================
[program:api-gateway]
command=uvicorn main:app --host 0.0.0.0 --port 8056 --log-level %(ENV_LOG_LEVEL)s
directory=/app/api-gateway
autostart=true
autorestart=true
priority=60
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",ADMIN_API_URL="http://localhost:8057",BOT_MANAGER_URL="http://localhost:8080",TRANSCRIPTION_COLLECTOR_URL="http://localhost:8123",MCP_URL="http://localhost:18888",REDIS_URL="%(ENV_REDIS_URL)s"

; =============================================================================
; MCP Service
; Model Context Protocol service for Claude/Cursor integration
; =============================================================================
[program:mcp]
command=uvicorn main:app --host 0.0.0.0 --port 18888 --log-level %(ENV_LOG_LEVEL)s
directory=/app/mcp
autostart=true
autorestart=true
priority=70
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",API_GATEWAY_URL="http://localhost:8056",LOG_LEVEL="%(ENV_LOG_LEVEL)s"

; =============================================================================
; TTS Service
; Text-to-speech via OpenAI API for voice agent / /speak command
; =============================================================================
[program:tts-service]
command=uvicorn main:app --host 0.0.0.0 --port 8059 --log-level %(ENV_LOG_LEVEL)s
directory=/app/tts-service
autostart=true
autorestart=true
priority=45
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s"

; =============================================================================
; Process Groups - for easier management
; =============================================================================
[group:vexa-core]
programs=redis,xvfb,pulseaudio,pulseaudio-setup,whisperlive,transcription-collector,admin-api,bot-manager,api-gateway,mcp,tts-service
priority=999
